# nginx.supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/local/openresty/bin/openresty -g "daemon off;"
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE

# Watch for changes in Nginx conf directories and restart Nginx when a config change occured
[program:nginx-reload]
command=bash -c 'while inotifywait -q -r -e create,delete,modify,move,attrib --exclude "/\." /etc/nginx/conf.d/ /usr/local/openresty/nginx/conf/; do nginx -t && nginx -s reload; done'
